# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
childlogdir=/var/log/supervisor
strip_ansi=true

[program:server]
environment=HOST_NAME="http://127.0.0.1"
command=node server.js
stdout_logfile=/tmp/web.log
stderr_logfile=/tmp/web.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[inet_http_server]
port = *:9001
username = dev 
password = dev
